# backend/warungkita.service
[Unit]
Description=WarungKita API Backend
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/warungkita/backend
ExecStart=/home/ubuntu/.hermes/hermes-agent/venv/bin/python -m uvicorn main:app --host 127.0.0.1 --port 8000 --log-level info
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/ubuntu/.hermes/hermes-agent/venv/lib/python3.11/site-packages
EnvironmentFile=-/home/ubuntu/.hermes/.env
EnvironmentFile=-/home/ubuntu/warungkita/backend/.env.admin

[Install]
WantedBy=multi-user.target
